# x-raider.service
[Unit]
Description=X-Raider Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/path/to/your/x-raider
ExecStart=/usr/bin/node telegramBot-cli.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=NODE_PATH=/path/to/your/x-raider/node_modules

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/path/to/your/x-raider/logs /path/to/your/x-raider/data /path/to/your/x-raider/temp

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
